# supervisor.conf
# Supervisor Configuration for Cargo Viewer Production
#
# GOAL: Configure Supervisor to manage Django application processes
#
# Features:
# - Gunicorn web server
# - Celery worker
# - Celery beat scheduler
# - Automatic restart on failure
# - Logging configuration
# - Process monitoring

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
chown=root:root

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# Gunicorn Web Server
[program:cargo-viewer-web]
command=/opt/cargo-viewer/venv/bin/gunicorn config.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 4 \
    --threads 4 \
    --worker-class sync \
    --worker-connections 1000 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --timeout 120 \
    --keepalive 5 \
    --log-level info \
    --access-logfile - \
    --error-logfile - \
    --capture-output \
    --enable-stdio-inheritance \
    --pid /tmp/gunicorn.pid
directory=/opt/cargo-viewer
user=deploy
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
killasgroup=true
stopwaitsecs=60
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cargo-viewer-web-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/cargo-viewer-web-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=
    DJANGO_SETTINGS_MODULE="config.settings.production",
    PATH="/opt/cargo-viewer/venv/bin:%(ENV_PATH)s"

# Celery Worker - Default Queue
[program:cargo-viewer-celery]
command=/opt/cargo-viewer/venv/bin/celery -A config worker \
    --loglevel=info \
    --concurrency=4 \
    --max-tasks-per-child=1000 \
    --time-limit=300 \
    --soft-time-limit=240 \
    --prefetch-multiplier=4 \
    -Q default,cargos,integrations,notifications
directory=/opt/cargo-viewer
user=deploy
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
killasgroup=true
stopwaitsecs=60
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cargo-viewer-celery-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/cargo-viewer-celery-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=
    DJANGO_SETTINGS_MODULE="config.settings.production",
    PATH="/opt/cargo-viewer/venv/bin:%(ENV_PATH)s",
    C_FORCE_ROOT="1"

# Celery Worker - High Priority Queue
[program:cargo-viewer-celery-high]
command=/opt/cargo-viewer/venv/bin/celery -A config worker \
    --loglevel=info \
    --concurrency=2 \
    --max-tasks-per-child=500 \
    --time-limit=180 \
    --soft-time-limit=150 \
    --prefetch-multiplier=2 \
    -Q high_priority
directory=/opt/cargo-viewer
user=deploy
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
killasgroup=true
stopwaitsecs=60
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cargo-viewer-celery-high-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/cargo-viewer-celery-high-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=
    DJANGO_SETTINGS_MODULE="config.settings.production",
    PATH="/opt/cargo-viewer/venv/bin:%(ENV_PATH)s",
    C_FORCE_ROOT="1"

# Celery Beat Scheduler
[program:cargo-viewer-celerybeat]
command=/opt/cargo-viewer/venv/bin/celery -A config beat \
    --loglevel=info \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/opt/cargo-viewer
user=deploy
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
killasgroup=true
stopwaitsecs=60
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cargo-viewer-celerybeat-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/cargo-viewer-celerybeat-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=
    DJANGO_SETTINGS_MODULE="config.settings.production",
    PATH="/opt/cargo-viewer/venv/bin:%(ENV_PATH)s",
    C_FORCE_ROOT="1"

# Celery Flower Monitoring (Optional)
[program:cargo-viewer-flower]
command=/opt/cargo-viewer/venv/bin/celery -A config flower \
    --port=5555 \
    --broker=redis://127.0.0.1:6379/1 \
    --basic_auth=admin:%(ENV_FLOWER_PASSWORD)s
directory=/opt/cargo-viewer
user=deploy
numprocs=1
autostart=false
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
killasgroup=true
stopwaitsecs=60
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cargo-viewer-flower-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/cargo-viewer-flower-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=
    DJANGO_SETTINGS_MODULE="config.settings.production",
    PATH="/opt/cargo-viewer/venv/bin:%(ENV_PATH)s",
    C_FORCE_ROOT="1"

# Group: All Cargo Viewer Processes
[group:cargo-viewer]
programs=cargo-viewer-web,cargo-viewer-celery,cargo-viewer-celery-high,cargo-viewer-celerybeat

# Event Listener for Process Monitoring
[eventlistener:memmon]
command=/usr/local/bin/memmon -a 200MB -m deploy@example.com
events=TICK_60

# Event Listener for Process Crash Notification
[eventlistener:crashmail]
command=/usr/local/bin/crashmail -a -m deploy@example.com
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL

# Event Listener for Process State Changes
[eventlistener:superlance]
command=/usr/local/bin/httpok -p /health/ http://127.0.0.1:8000/ cargo-viewer-web
events=TICK_60

# Include additional configurations
[include]
files = /etc/supervisor/conf.d/*.conf
